# claude-resident-health@.service
[Unit]
Description=Claude Resident 헬스체크 (%i)

[Service]
Type=oneshot
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=%h/.local/bin/claude-resident %i health
