# claude-resident-restart@.service
[Unit]
Description=Claude Resident 일일 재시작 실행 (%i)
# timer에 의해 호출되는 oneshot 서비스

[Service]
Type=oneshot
# main service의 ExecStop/ExecStart를 통해 tmux 세션이 claude-resident@%i.service
# cgroup 아래에서 일관되게 관리되도록 재시작한다.
ExecStart=/usr/bin/systemctl --user restart claude-resident@%i.service
