# claude-resident-health@.service
[Unit]
Description=Claude Resident 헬스체크 (%i)

[Service]
Type=oneshot
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash -lc 'systemctl --user is-active --quiet claude-resident@%i.service || exit 0; tmux has-session -t claude-resident-%i 2>/dev/null || systemctl --user restart claude-resident@%i.service'
