# claude-resident-health@.timer
[Unit]
Description=Claude Resident 헬스체크 타이머 (%i)

[Timer]
OnBootSec=3min
OnUnitActiveSec=3min
Persistent=true

[Install]
WantedBy=timers.target
