# claude-resident-restart@.service
[Unit]
Description=Claude Resident 일일 재시작 실행 (%i)
# timer에 의해 호출되는 oneshot 서비스

[Service]
Type=oneshot
# shutdown 시퀀스 실행 (recent.md 저장 시도) 후 재시작
ExecStart=/bin/bash -c '/home/chjee/.local/bin/claude-resident %i shutdown && sleep 3 && /home/chjee/.local/bin/claude-resident %i start'
# shutdown이 실패해도(recent.md 저장 실패) 재시작은 진행
